# profile.d/which2.sh
# Initialization script for bash and sh

alias which='(alias; declare -f) | /usr/bin/which --tty-only --read-alias --read-functions --show-tilde --show-dot'
